# deploy/systemd/reco-demo.service
[Unit]
Description=Reco demo (Dashdoor + engine sidecars + agent)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/reco-demo
EnvironmentFile=-/etc/reco-demo/env
ExecStart=/usr/bin/docker compose -f config/docker-compose.demo.yaml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f config/docker-compose.demo.yaml down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
